# deploy/systemd/eibrain-vision-hailo.service
[Unit]
Description=eibrain Hailo vision state service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/home/darrow/eibrain
EnvironmentFile=-/home/darrow/eibrain/.env
ExecStart=/home/darrow/eibrain/.venv/bin/python -m apps.body_runtime.vision_hailo_service --config /home/darrow/eibrain/config/eibrain.honjia.local.yaml --backend gstreamer --interval-s 0.1 --security-vision-always-on
Restart=always
RestartSec=3
StandardOutput=append:/home/darrow/eibrain/.tmp-test-artifacts/eibrain-vision-hailo.service.log
StandardError=append:/home/darrow/eibrain/.tmp-test-artifacts/eibrain-vision-hailo.service.log

[Install]
WantedBy=default.target
